# root/etc/s6-overlay/s6-rc.d/p12-sshfs-setup/run
#!/usr/bin/with-contenv bash
set -e

# bring your own custom sshfs script
echo "Nothing to do for SSHFS.";
