# root/etc/s6-overlay/s6-rc.d/p12-sshfs-setup/run
#!/usr/bin/with-contenv bash
set -e

vecho () { if [ "${S6_VERBOSITY:-1}" -gt 0 ]; then echo "[$0] $@"; fi; }

# bring your own custom sshfs script
vecho "Nothing to do for SSHFS.";
