# root/etc/s6-overlay/s6-rc.d/p12-sshfs-setup/up
/etc/s6-overlay/s6-rc.d/p12-sshfs-setup/run
