# root/etc/services.d/ssh/run
#!/usr/bin/with-contenv bash

/usr/sbin/sshd -Dd -e4 -f /etc/ssh/sshd_config;

