# root/etc/s6-overlay/s6-rc.d/p11-autossh-setup/run
#!/usr/bin/with-contenv bash
set -e

# bring your own custom autossh script
echo "Nothing to do for AutoSSH.";
